# jvb/rootfs/etc/supervisor/conf.d/10_jvb.conf
[program:jvb]
command=/bin/bash -c "/usr/bin/service_wrapper_jvb"
autorestart=unexpected
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true